# ini/inittab
::sysinit:clear
::sysinit:cat /etc/logo.txt
::restart:/sbin/init
::shutdown:sync
::shutdown:umount -a
::ctrlaltdel:/sbin/reboot
::respawn:/bin/cttyhack /etc/shell.sh
