# root/etc/inittab
::sysinit:clear
::sysinit:echo "  Monkey See, Monkey Do LINUX"
::restart:/sbin/init
::shutdown:sync
::ctrlaltdel:/sbin/reboot
::respawn:/bin/cttyhack /etc/shell.sh
