# root/etc/inittab
::sysinit:clear
::restart:/sbin/init
::shutdown:sync
::ctrlaltdel:/sbin/reboot
::respawn:/bin/cttyhack login -f msmd
